# services.d/forwarder-ui/octoblu-forwarder-ui-sidekick@.service
[Unit]
Description=octoblu-forwarder-ui-sidekick
After=octoblu-forwarder-ui@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-forwarder-ui-sidekick-%i \
  && /usr/bin/etcdctl get /kronk/wrapper/run > /run/kronk/octoblu-forwarder-ui-sidekick-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-forwarder-ui-sidekick-%i/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-forwarder-ui-sidekick-%i/kronk-run \
    /kronk/sidekick/start octoblu-forwarder-ui-%i \
    octoblu \
    forwarder-ui \
    %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-forwarder-ui@%i.service
