# services.d/minor-server-update/octoblu-minor-server-update-sidekick@.service
[Unit]
Description=octoblu-minor-server-update-sidekick
After=octoblu-minor-server-update@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-minor-server-update-sidekick-%i \
  && /usr/bin/etcdctl get /kronk/wrapper/run > /run/kronk/octoblu-minor-server-update-sidekick-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-minor-server-update-sidekick-%i/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-minor-server-update-sidekick-%i/kronk-run \
    /kronk/sidekick/start octoblu-minor-server-update-%i \
    octoblu \
    minor-server-update \
    %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-minor-server-update@%i.service
