# services.d/refresh-token-worker/octoblu-refresh-token-worker@.service
[Unit]
Description=octoblu-refresh-token-worker
Wants=octoblu-refresh-token-worker-register@%i.service
Wants=octoblu-refresh-token-worker-sidekick@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=300
StartLimitInterval=30
StartLimitBurst=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-refresh-token-worker-%i \
  && /usr/bin/etcdctl get /kronk/wrapper/run > /run/kronk/octoblu-refresh-token-worker-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-refresh-token-worker-%i/kronk-run \
"
ExecStartPre=/bin/bash -c " \
  /run/kronk/octoblu-refresh-token-worker-%i/kronk-run \
    /kronk/service/start-pre octoblu-refresh-token-worker-%i \
      octoblu \
      refresh-token-worker \
      %i \
"
ExecStart=/bin/bash -c " \
  /run/kronk/octoblu-refresh-token-worker-%i/kronk-run \
    /kronk/service/start octoblu-refresh-token-worker-%i \
      octoblu \
      refresh-token-worker \
      %i \
"
ExecStop=/bin/bash -c " \
  /run/kronk/octoblu-refresh-token-worker-%i/kronk-run \
    /kronk/service/stop octoblu-refresh-token-worker-%i \
      octoblu \
      refresh-token-worker \
      %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
MachineMetadata=Services=true
