# services.d/meshblu-authenticator-facebook/octoblu-meshblu-authenticator-facebook-sidekick@.service
[Unit]
Description=octoblu-meshblu-authenticator-facebook-sidekick
After=octoblu-meshblu-authenticator-facebook@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-meshblu-authenticator-facebook-sidekick-%i \
  && /usr/bin/etcdctl get /kronk/wrapper/run > /run/kronk/octoblu-meshblu-authenticator-facebook-sidekick-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-meshblu-authenticator-facebook-sidekick-%i/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-meshblu-authenticator-facebook-sidekick-%i/kronk-run \
    /kronk/sidekick/start octoblu-meshblu-authenticator-facebook-%i \
    octoblu \
    meshblu-authenticator-facebook \
    %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-meshblu-authenticator-facebook@%i.service
