# services.d/meshblu-authenticator-pin/octoblu-meshblu-authenticator-pin-register@.service
[Unit]
Description=octoblu-meshblu-authenticator-pin-register


[Service]
TimeoutStartSec=300
EnvironmentFile=/etc/environment
Restart=always
ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk to /tmp' \
  && /usr/bin/etcdctl get /kronk/register/start-pre > /tmp/octoblu-meshblu-authenticator-pin-register-%i-kronk-start-pre \
  && /usr/bin/etcdctl get /kronk/register/start > /tmp/octoblu-meshblu-authenticator-pin-register-%i-kronk-start \
  && /usr/bin/etcdctl get /kronk/register/stop > /tmp/octoblu-meshblu-authenticator-pin-register-%i-kronk-stop \
  && chmod +x /tmp/octoblu-meshblu-authenticator-pin-register-%i-kronk-* \
"
ExecStartPre=/bin/bash -c " \
  /tmp/octoblu-meshblu-authenticator-pin-register-%i-kronk-start-pre \
    octoblu-meshblu-authenticator-pin \
    %i \
    /octoblu/meshblu-authenticator-pin/register-env \
"
ExecStart=/bin/sh -c " \
  /tmp/octoblu-meshblu-authenticator-pin-register-%i-kronk-start \
    octoblu-meshblu-authenticator-pin \
    %i \
"
ExecStop=/bin/sh -c " \
  /tmp/octoblu-meshblu-authenticator-pin-register-%i-kronk-stop \
    octoblu-meshblu-authenticator-pin \
    %i \
"

[X-Fleet]
X-ConditionMachineOf=octoblu-meshblu-authenticator-pin@%i.service
