# services.d/slurry-twitter/octoblu-slurry-twitter-sidekick@.service
[Unit]
Description=octoblu-slurry-twitter-sidekick
After=octoblu-slurry-twitter@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-slurry-twitter-sidekick-%i \
  && /usr/bin/etcdctl get /kronk/wrapper/run > /run/kronk/octoblu-slurry-twitter-sidekick-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-slurry-twitter-sidekick-%i/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-slurry-twitter-sidekick-%i/kronk-run \
    /kronk/sidekick/start octoblu-slurry-twitter-%i \
    octoblu \
    slurry-twitter \
    %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-slurry-twitter@%i.service
