# services.d/data-forwarder-azure-service-bus/octoblu-data-forwarder-azure-service-bus@.service
[Unit]
Description=octoblu-data-forwarder-azure-service-bus
Wants=octoblu-data-forwarder-azure-service-bus-register@%i.service
Wants=octoblu-data-forwarder-azure-service-bus-sidekick@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=300
StartLimitInterval=30
StartLimitBurst=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-data-forwarder-azure-service-bus-%i \
  && /usr/bin/etcdctl get /kronk/wrapper/run > /run/kronk/octoblu-data-forwarder-azure-service-bus-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-data-forwarder-azure-service-bus-%i/kronk-run \
"
ExecStartPre=/bin/bash -c " \
  /run/kronk/octoblu-data-forwarder-azure-service-bus-%i/kronk-run \
    /kronk/service/start-pre octoblu-data-forwarder-azure-service-bus-%i \
      octoblu \
      data-forwarder-azure-service-bus \
      %i \
"
ExecStart=/bin/bash -c " \
  /run/kronk/octoblu-data-forwarder-azure-service-bus-%i/kronk-run \
    /kronk/service/start octoblu-data-forwarder-azure-service-bus-%i \
      octoblu \
      data-forwarder-azure-service-bus \
      %i \
"
ExecStop=/bin/bash -c " \
  /run/kronk/octoblu-data-forwarder-azure-service-bus-%i/kronk-run \
    /kronk/service/stop octoblu-data-forwarder-azure-service-bus-%i \
      octoblu \
      data-forwarder-azure-service-bus \
      %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
MachineMetadata=Services=true
