# services.d/nanocyte-engine-worker/octoblu-nanocyte-engine-worker@.service
[Unit]
Description=octoblu-nanocyte-engine-worker
Wants=octoblu-nanocyte-engine-worker-register@%i.service
Wants=octoblu-nanocyte-engine-worker-sidekick@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=300
StartLimitInterval=30
StartLimitBurst=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-nanocyte-engine-worker-%i \
  && /usr/bin/curl -s http://kronk.octoblu.com/latest/wrapper/run > /run/kronk/octoblu-nanocyte-engine-worker-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-nanocyte-engine-worker-%i/kronk-run \
"
ExecStartPre=/bin/bash -c " \
  /run/kronk/octoblu-nanocyte-engine-worker-%i/kronk-run \
    /kronk/service/start-pre octoblu-nanocyte-engine-worker-%i \
      octoblu \
      nanocyte-engine-worker \
      %i \
"
ExecStart=/bin/bash -c " \
  /run/kronk/octoblu-nanocyte-engine-worker-%i/kronk-run \
    /kronk/service/start octoblu-nanocyte-engine-worker-%i \
      octoblu \
      nanocyte-engine-worker \
      %i \
"
ExecStop=/bin/bash -c " \
  /run/kronk/octoblu-nanocyte-engine-worker-%i/kronk-run \
    /kronk/service/stop octoblu-nanocyte-engine-worker-%i \
      octoblu \
      nanocyte-engine-worker \
      %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
MachineMetadata=Services=true
