# services.d/meshblu-request-queue-length-logger/meshblu-request-queue-length-logger.timer
[Unit]
Description=Run and log the length of the meshblu-request queue

[Timer]
OnBootSec=3 minutes
OnUnitActiveSec=1 minute

[Install]
WantedBy=timers.target

[X-Fleet]
X-ConditionMachineOf=meshblu-request-queue-length-logger.service
