# services.d/slurry-salesforce/octoblu-slurry-salesforce@.service
[Unit]
Description=octoblu-slurry-salesforce
Wants=octoblu-slurry-salesforce-register@%i.service
Wants=octoblu-slurry-salesforce-sidekick@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=300
StartLimitInterval=30
StartLimitBurst=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-slurry-salesforce-%i \
  && /usr/bin/curl -s http://kronk.octoblu.com/latest/wrapper/run > /run/kronk/octoblu-slurry-salesforce-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-slurry-salesforce-%i/kronk-run \
"
ExecStartPre=/bin/bash -c " \
  /run/kronk/octoblu-slurry-salesforce-%i/kronk-run \
    /kronk/service/start-pre octoblu-slurry-salesforce-%i \
      octoblu \
      slurry-salesforce \
      %i \
"
ExecStart=/bin/bash -c " \
  /run/kronk/octoblu-slurry-salesforce-%i/kronk-run \
    /kronk/service/start octoblu-slurry-salesforce-%i \
      octoblu \
      slurry-salesforce \
      %i \
"
ExecStop=/bin/bash -c " \
  /run/kronk/octoblu-slurry-salesforce-%i/kronk-run \
    /kronk/service/stop octoblu-slurry-salesforce-%i \
      octoblu \
      slurry-salesforce \
      %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
MachineMetadata=Services=true
