# services.d/slurry-salesforce/octoblu-slurry-salesforce-sidekick@.service
[Unit]
Description=octoblu-slurry-salesforce-sidekick
After=octoblu-slurry-salesforce@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-slurry-salesforce-sidekick-%i \
  && /usr/bin/curl -s http://kronk.octoblu.com/latest/wrapper/run > /run/kronk/octoblu-slurry-salesforce-sidekick-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-slurry-salesforce-sidekick-%i/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-slurry-salesforce-sidekick-%i/kronk-run \
    /kronk/sidekick/start octoblu-slurry-salesforce-%i \
    octoblu \
    slurry-salesforce \
    %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-slurry-salesforce@%i.service
