# services.d/governator-service/octoblu-governator-service-sidekick@.service
[Unit]
Description=octoblu-governator-service-sidekick
After=octoblu-governator-service@%i.service

[Service]
TimeoutStartSec=30
Restart=always


ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk to /tmp' \
  && /usr/bin/etcdctl get /octoblu/kronk/sidekick/start > /tmp/octoblu-governator-service-sidekick-%i-kronk-start \
  && chmod +x /tmp/octoblu-governator-service-sidekick-%i-kronk-* \
"
ExecStart=/usr/bin/etcdctl exec-watch \
  /octoblu/governator-service/restart \
  -- /tmp/octoblu-governator-service-sidekick-%i-kronk-start 10 /octoblu/governator-service/docker_url octoblu-governator-service-%i octoblu-governator-service-register-%i

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-governator-service@%i.service
