# services.d/service-state-service/octoblu-service-state-service-sidekick@.service
[Unit]
Description=octoblu-service-state-service-sidekick
After=octoblu-service-state-service@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-service-state-service-sidekick-%i \
  && /usr/bin/curl -s http://kronk.octoblu.com/latest/wrapper/run > /run/kronk/octoblu-service-state-service-sidekick-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-service-state-service-sidekick-%i/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-service-state-service-sidekick-%i/kronk-run \
    /kronk/sidekick/start octoblu-service-state-service-%i \
    octoblu \
    service-state-service \
    %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-service-state-service@%i.service
