# services.d/aws-iam/octoblu-aws-iam-sidekick@.service
[Unit]
Description=octoblu-aws-iam-sidekick
After=octoblu-aws-iam@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-aws-iam-sidekick-%i \
  && /usr/bin/curl -s http://kronk.octoblu.com/latest/wrapper/run > /run/kronk/octoblu-aws-iam-sidekick-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-aws-iam-sidekick-%i/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-aws-iam-sidekick-%i/kronk-run \
    /kronk/sidekick/start octoblu-aws-iam-%i \
    octoblu \
    aws-iam \
    %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-aws-iam@%i.service
