# services.d/aws-iam/octoblu-aws-iam@.service
[Unit]
Description=octoblu-aws-iam
Wants=octoblu-aws-iam-register@%i.service
Wants=octoblu-aws-iam-sidekick@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=300
StartLimitInterval=30
StartLimitBurst=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-aws-iam-%i \
  && /usr/bin/etcdctl get /kronk/wrapper/run > /run/kronk/octoblu-aws-iam-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-aws-iam-%i/kronk-run \
"
ExecStartPre=/bin/bash -c " \
  /run/kronk/octoblu-aws-iam-%i/kronk-run \
    /kronk/service/start-pre octoblu-aws-iam-%i \
      octoblu \
      aws-iam \
      %i \
"
ExecStart=/bin/bash -c " \
  /run/kronk/octoblu-aws-iam-%i/kronk-run \
    /kronk/service/start octoblu-aws-iam-%i \
      octoblu \
      aws-iam \
      %i \
"
ExecStop=/bin/bash -c " \
  /run/kronk/octoblu-aws-iam-%i/kronk-run \
    /kronk/service/stop octoblu-aws-iam-%i \
      octoblu \
      aws-iam \
      %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
MachineMetadata=Services=true
