# services.d/meshblu-balancer-socket.io/octoblu-meshblu-balancer-socket.io-sidekick.service
[Unit]
Description=octoblu-meshblu-balancer-socket.io-sidekick
After=octoblu-meshblu-balancer-socket.io.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-meshblu-balancer-socket.io-sidekick \
  && /usr/bin/etcdctl get /kronk/wrapper/run > /run/kronk/octoblu-meshblu-balancer-socket.io-sidekick/kronk-run \
  && chmod +x /run/kronk/octoblu-meshblu-balancer-socket.io-sidekick/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-meshblu-balancer-socket.io-sidekick/kronk-run \
    /kronk/sidekick/start octoblu-meshblu-balancer-socket.io \
    octoblu \
    meshblu-balancer-socket.io \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-meshblu-balancer-socket.io.service
