# services.d/meshblu-verifier-xmpp/octoblu-meshblu-verifier-xmpp.timer
[Unit]
Description=meshblu-verifier-xmpp timer

[Timer]
OnBootSec=3 minutes
OnUnitActiveSec=1 minute

[Install]
WantedBy=timers.target

[X-Fleet]
MachineMetadata=Services=true
X-ConditionMachineOf=octoblu-meshblu-verifier-xmpp.service
