# services.d/data-forwarder-mongodb/octoblu-data-forwarder-mongodb-sidekick@.service
[Unit]
Description=octoblu-data-forwarder-mongodb-sidekick
After=octoblu-data-forwarder-mongodb@%i.service

[Service]
TimeoutStartSec=30
Restart=always


ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-data-forwarder-mongodb-sidekick-%i \
  && /usr/bin/etcdctl get /kronk/wrapper/run > /run/kronk/octoblu-data-forwarder-mongodb-sidekick-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-data-forwarder-mongodb-sidekick-%i/kronk-run \
"
ExecStart=/usr/bin/etcdctl exec-watch \
  /octoblu/data-forwarder-mongodb/restart \
  -- /run/kronk/octoblu-data-forwarder-mongodb-sidekick-%i/kronk-run \
    /kronk/sidekick/start octoblu-data-forwarder-mongodb-sidekick-%i \
    10 /octoblu/data-forwarder-mongodb/docker_url octoblu-data-forwarder-mongodb-%i octoblu-data-forwarder-mongodb-register-%i

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-data-forwarder-mongodb@%i.service
