# services.d/nanocyte-interval-redis/octoblu-nanocyte-interval-redis-register@.service
[Unit]
Description=octoblu-nanocyte-interval-redis-register


[Service]
TimeoutStartSec=300
EnvironmentFile=/etc/environment
Restart=always
ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk to /tmp' \
  && /usr/bin/etcdctl get /kronk/register/start-pre > /tmp/octoblu-nanocyte-interval-redis-register-%i-kronk-start-pre \
  && /usr/bin/etcdctl get /kronk/register/start > /tmp/octoblu-nanocyte-interval-redis-register-%i-kronk-start \
  && /usr/bin/etcdctl get /kronk/register/stop > /tmp/octoblu-nanocyte-interval-redis-register-%i-kronk-stop \
  && chmod +x /tmp/octoblu-nanocyte-interval-redis-register-%i-kronk-* \
"
ExecStartPre=/bin/bash -c " \
  /tmp/octoblu-nanocyte-interval-redis-register-%i-kronk-start-pre \
    octoblu-nanocyte-interval-redis \
    %i \
    /octoblu/nanocyte-interval-redis/register-env \
"
ExecStart=/bin/sh -c " \
  /tmp/octoblu-nanocyte-interval-redis-register-%i-kronk-start \
    octoblu-nanocyte-interval-redis \
    %i \
"
ExecStop=/bin/sh -c " \
  /tmp/octoblu-nanocyte-interval-redis-register-%i-kronk-stop \
    octoblu-nanocyte-interval-redis \
    %i \
"

[X-Fleet]
X-ConditionMachineOf=octoblu-nanocyte-interval-redis@%i.service
