# services.d/endo-sendgrid/octoblu-endo-sendgrid-sidekick@.service
[Unit]
Description=octoblu-endo-sendgrid-sidekick
After=octoblu-endo-sendgrid@%i.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-endo-sendgrid-sidekick-%i \
  && /usr/bin/curl -s http://kronk.octoblu.com/latest/wrapper/run > /run/kronk/octoblu-endo-sendgrid-sidekick-%i/kronk-run \
  && chmod +x /run/kronk/octoblu-endo-sendgrid-sidekick-%i/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-endo-sendgrid-sidekick-%i/kronk-run \
    /kronk/sidekick/start octoblu-endo-sendgrid-%i \
    octoblu \
    endo-sendgrid \
    %i \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-endo-sendgrid@%i.service
