# services.d/meshblu-amqp/octoblu-meshblu-amqp-sidekick.service
[Unit]
Description=octoblu-meshblu-amqp-sidekick

[Service]
TimeoutStartSec=30
Restart=always
ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk to /tmp' \
  && /usr/bin/etcdctl get /kronk/sidekick/start > /tmp/octoblu-meshblu-amqp-sidekick-kronk-start \
  && chmod +x /tmp/octoblu-meshblu-amqp-sidekick-kronk-* \
"
ExecStart=/usr/bin/etcdctl exec-watch \
  /octoblu/meshblu-amqp/restart \
  -- /tmp/octoblu-meshblu-amqp-sidekick-kronk-start 10 /octoblu/meshblu-amqp/docker_url octoblu-meshblu-amqp octoblu-meshblu-amqp-register

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-meshblu-amqp.service
