# services.d/meshblu-balancer-xmpp/octoblu-meshblu-balancer-xmpp-register-sidekick.service
[Unit]
Description=octoblu-meshblu-balancer-xmpp-sidekick
After=octoblu-meshblu-balancer-xmpp.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=30
Restart=always

ExecStartPre=/bin/bash -c " \
  echo 'Writing kronk' \
  && mkdir -p /run/kronk/octoblu-meshblu-balancer-xmpp-sidekick \
  && /usr/bin/curl -s http://kronk.octoblu.com/latest/wrapper/run > /run/kronk/octoblu-meshblu-balancer-xmpp-sidekick/kronk-run \
  && chmod +x /run/kronk/octoblu-meshblu-balancer-xmpp-sidekick/kronk-run \
"
ExecStart=/bin/sh -c " \
  /run/kronk/octoblu-meshblu-balancer-xmpp-sidekick/kronk-run \
    /kronk/balancer-register-sidekick/start octoblu-meshblu-balancer-xmpp \
    octoblu \
    meshblu-balancer-xmpp \
    meshblu-core-protocol-adapter-xmpp \
"

[Install]
WantedBy=multi-user.target

[X-Fleet]
X-ConditionMachineOf=octoblu-meshblu-balancer-xmpp.service
